# iso/rootfs/etc/xdg/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=Install Tejas Linux
Exec=pkexec calamares
X-GNOME-Autostart-enabled=true
